# systemd/squadjs-restart.timer
[Unit]
Description=Restart SquadJS daily

[Timer]
OnCalendar=06:10

[Install]
WantedBy=timers.target
